# roles/prometheus_agent/templates/node_exporter.service.j2
[Unit]
Description=NodeExporter

[Service]
User=root
TimeoutStartSec=0
ExecStart={{ node_exporter_bin_path }} --web.listen-address=:{{ node_exporter_port }}

[Install]
WantedBy=multi-user.target
